# .config/autostart/ssh-add.desktop
[Desktop Entry]
Exec=ssh-add -q ~/.ssh/id_ed25519 < /dev/null
Name=ssh-add
Type=Application
